# filebeat/rootfs/etc/services.d/filebeat/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Filebeat
# Runs the filebeat daemon
# ==============================================================================
bashio::log.info "Starting filebeat..."
exec ls / && find /share && tail -f /share/hassio/homeassistant/homeassistant.log
#exec filebeat -c /etc/filebeat/filebeat.yml
