# filebeat/rootfs/etc/services.d/filebeat/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Filebeat
# Runs the filebeat daemon
# ==============================================================================
bashio::log.info "Starting filebeat..."
exec find /var
#exec filebeat -c /etc/filebeat/filebeat.yml
